# freebsd/scripts/freebsd/pkg-plist
etc/martin.d/martin.env.sample
etc/rc.conf.d/martin
etc/rc.d/martin
libexec/martin/martin
libexec/martin/martin-cp
